# deploy/aigoproxy.service
# aigoproxy.service — systemd unit (system mode, runs as aigoproxy user)
#
# Why system mode: user-mode systemd can't grant CAP_NET_BIND_SERVICE
# (needed for :80, :443) without CAP_SETFCAP, which most distros don't
# grant to the install user. System mode + dedicated user is the
# standard production pattern for a public reverse proxy.
#
# Install: see ../install.sh (creates user + dirs + setcap on the binary)
# Status:  systemctl status aigoproxy
# Logs:    journalctl -u aigoproxy -f
# Reload:  systemctl reload aigoproxy          # sends SIGHUP
# Stop:    systemctl stop aigoproxy

[Unit]
Description=aigoproxy — Tailscale-aware reverse proxy
Documentation=https://github.com/biodoia/aigoproxy
After=network-online.target tailscaled.service
Wants=network-online.target

[Service]
Type=simple
User=aigoproxy
Group=aigoproxy
# When run as the aigoproxy user, we use setcap to bind :80.
# setcap is set by install.sh; if it's missing, bind to a high port
# with `-addr :8080`.
ExecStart=/usr/local/bin/aigoproxy -addr :80 -config /etc/aigoproxy/config.yaml -data /var/lib/aigoproxy
# Graceful config reload: sends SIGHUP, which aigoproxy handles by
# re-reading config.yaml and rebuilding the route table.
ExecReload=/bin/kill -HUP $MAINPID
# Stop gracefully — give the HTTP server 10s to drain in-flight requests.
TimeoutStopSec=10s
KillMode=mixed
Restart=on-failure
RestartSec=5s

# Hardening. We do NOT set NoNewPrivileges=true because the binary uses
# file capabilities (setcap cap_net_bind_service) to bind :80. File
# capabilities act like a setuid bit, which NoNewPrivileges blocks.
# Instead, we use AmbientCapabilities to grant the cap directly to the
# process and to its children, without needing root.
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/aigoproxy /var/lib/tailscale /home/%i/.local/share/tailscale
ReadOnlyPaths=/etc/aigoproxy
PrivateTmp=true
PrivateDevices=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6 AF_NETLINK
RestrictNamespaces=true
LockPersonality=true

# Grant CAP_NET_BIND_SERVICE so the process can bind :80 and :443
# without running as root. This works alongside the setcap applied
# to the binary at install time.
AmbientCapabilities=CAP_NET_BIND_SERVICE
CapabilityBoundingSet=CAP_NET_BIND_SERVICE

# We bind :80 which is a privileged port. setcap is applied at install
# time (see install.sh), so we don't need root or AmbientCapabilities.

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=aigoproxy

[Install]
WantedBy=multi-user.target
